# scripts/services/ai-studio.timer
[Unit]
Description=Timer for AI Studio Automation Service

[Timer]
OnBootSec=30s
Unit=ai-studio.service

[Install]
WantedBy=timers.target